# federvise.service
[Unit]
Description=Federvise ActivityPub Blog
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/federvise
ExecStart=/usr/bin/node dist/index.js
Restart=on-failure
RestartSec=10

# 环境变量
EnvironmentFile=/var/www/federvise/.env

# 安全设置
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/federvise/data

[Install]
WantedBy=multi-user.target
